# backend/dist/unit/corrections.service
[Unit]
Description=MineMonitor Corrections
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/python3.6 /usr/local/bin/corrections

[Install]
WantedBy=multi-user.target