# backend/dist/unit/overview.service
[Unit]
Description=MineMonitor Overview
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/python3.6 /usr/local/bin/overview

[Install]
WantedBy=multi-user.target