# backend/dist/unit/services.service
[Unit]
Description=MineMonitor Services
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/python3.6 /usr/local/bin/services

[Install]
WantedBy=multi-user.target